# root/etc/systemd/system/update_fetcher.service
[Unit]
Description=Autodownload updates without installing

[Service]
Type=oneshot
ExecStart=pacman -Syuw --noconfirm
